# .config/systemd/user/i3-graphical-session.service
[Unit]
Description=trigger graphical-session
Wants=graphical-session.target

[Service]
ExecStart=/bin/false
